# x2_api/x2-mqtt-client.service
[Unit]
Description=X2 MQTT Client (AWS IoT Core)
After=agibot_software.service
Wants=agibot_software.service

[Service]
Type=simple
User=run
WorkingDirectory=/home/run/x2_api
ExecStart=/home/run/x2_api/mqtt_run.sh
Restart=on-failure
RestartSec=5
Environment=X2_ROBOT_ID=x2-001
Environment=X2_IOT_ENDPOINT=a1thbiemoccm90-ats.iot.us-east-2.amazonaws.com
Environment=X2_CERT_PATH=/home/run/x2_api/certs

[Install]
WantedBy=multi-user.target
